# pi/config/systemd/user/register-device.service
[Unit]
Description=A job to register a device

[Service]
Type=simple
ExecStart=/bin/register.sh

[Install]
WantedBy=multi-user.target
